# sysconf/ntp.service.d/configuration-exists.conf
[Unit]
ConditionPathExists=/etc/ntp.conf
